# supervisor/backend.ini
[program:backend]
command=/home/bastion/backend
user=bastion
environment=HOME="/home/bastion",SECRET="Tz3O&Z*x"